# contrib/init/btchdd.conf
description "BitcoinHD Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env xdiskD_BIN="/usr/bin/xdiskd"
env xdiskD_USER="xdisk"
env xdiskD_GROUP="xdisk"
env xdiskD_PIDDIR="/var/run/xdiskd"
# upstart can't handle variables constructed with other variables
env xdiskD_PIDFILE="/var/run/xdiskd/xdiskd.pid"
env xdiskD_CONFIGFILE="/etc/xdisk/xdisk.conf"
env xdiskD_DATADIR="/var/lib/xdiskd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # xdiskd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$xdiskD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run xdiskd."
        echo "The setting must appear in $xdiskD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"BitcoinHD Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$xdiskD_PIDDIR"
    chmod 0755 "$xdiskD_PIDDIR"
    chown $xdiskD_USER:$xdiskD_GROUP "$xdiskD_PIDDIR"
    chown $xdiskD_USER:$xdiskD_GROUP "$xdiskD_CONFIGFILE"
    chmod 0660 "$xdiskD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$xdiskD_PIDFILE" \
    --chuid $xdiskD_USER:$xdiskD_GROUP \
    --exec "xdiskD_BIN" \
    -- \
    -pid="$xdiskD_PIDFILE" \
    -conf="$xdiskD_CONFIGFILE" \
    -datadir="$xdiskD_DATADIR" \
    -disablewallet \
    -daemon

